# weather_station/research/weather-correction.service
[Unit]
Description=Weather Station — daily correction model retrain
After=network.target

[Service]
Type=oneshot
User=fejot
WorkingDirectory=/home/fejot/diplom/weather_station
ExecStart=/home/fejot/diplom/weather_station/venv/bin/python3 research/calibration_cli.py train-correction --yes
StandardOutput=journal
StandardError=journal
